# birdhouse/services/birdhouse-capture.service
[Unit]
Description=Birdhouse Camera — Capture Loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rafablazer
WorkingDirectory=/home/rafablazer/birdhouse
ExecStart=/usr/bin/python3 /home/rafablazer/birdhouse/capture.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
